# system/install/autoupdate.service
[Unit]
Description=Auto-update check

[Service]
Type=oneshot
ExecStart=/home/pi/mirror/system/run.sh
