# system/install/autoupdate.service
[Unit]
Description=Auto-update check

[Service]
Type=oneshot
ExecStart=%h/mirror/system/run.sh
