# system/install/reboot.service
[Unit]
Description=Reboot the system

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl reboot
