# system/install/screenon.timer
[Unit]
Description=Turn the screen on according to a schedule

[Timer]
OnCalendar=*-*-* 05:45:00
Persistent=true

[Install]
WantedBy=timers.target
