# system/install/screenon.service
[Unit]
Description=Turn the screen on

[Service]
Type=oneshot
ExecStart=%h/mirror/system/install/screenon.sh
