# system/install/screenoff.timer
[Unit]
Description=Turn the screen off according to a schedule

[Timer]
OnCalendar=*-*-* 23:00:00
Persistent=true

[Install]
WantedBy=timers.target
