# system/install/screenoff.service
[Unit]
Description=Turn the screen off

[Service]
Type=oneshot
ExecStart=%h/mirror/system/install/screenoff.sh
